# configure-zram.service
[Unit]
After=network.target

[Service]
ExecStart=/opt/scripts/enable-zram.sh

[Install]
WantedBy=default.target

